# weather-server.service
[Unit]
Description=Weather Alert Server
After=network.target

[Service]
WorkingDirectory=/home/user/wheather-GO
ExecStart=/home/user/wheather-GO/server
Restart=always
RestartSec=5
EnvironmentFile=/home/user/wheather-GO/.env

[Install]
WantedBy=multi-user.target
